# tests/docker/supervisord-minimal.conf
[supervisord]
nodaemon=true
logfile=/var/log/marlise/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:jack-dummy]
command=/usr/local/bin/start-jack-dummy.sh
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/jack-dummy.log
stdout_logfile=/var/log/marlise/jack-dummy.log
environment=HOME="/home/marlise",USER="marlise"
priority=100

[program:mock-mod-host]
command=python3 /opt/marlise/tests/docker/mock_mod_host.py
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/mod-host.log
stdout_logfile=/var/log/marlise/mod-host.log
environment=HOME="/home/marlise",USER="marlise"
priority=200
depends_on=jack-dummy

[program:mock-modhost-bridge]  
command=python3 /opt/marlise/tests/docker/mock_modhost_bridge.py
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/modhost-bridge.log
stdout_logfile=/var/log/marlise/modhost-bridge.log
environment=HOME="/home/marlise",USER="marlise"
priority=300
depends_on=mock-mod-host