# tests/docker/supervisord-test.conf
[supervisord]
nodaemon=true
user=marlise
logfile=/var/log/marlise/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:jack-dummy]
command=/usr/local/bin/start-jack-dummy.sh
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/jack-dummy.log
stdout_logfile=/var/log/marlise/jack-dummy.log
environment=HOME="/home/marlise",USER="marlise"
priority=100

[program:mod-host]
command=/opt/marlise/audio-engine/mod-host/mod-host -p 5555 -f 5556 -v
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/mod-host.log
stdout_logfile=/var/log/marlise/mod-host.log
environment=HOME="/home/marlise",USER="marlise"
priority=200
depends_on=jack-dummy

[program:modhost-bridge]
command=/opt/marlise/audio-engine/modhost-bridge/build/modhost-bridge
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/modhost-bridge.log
stdout_logfile=/var/log/marlise/modhost-bridge.log
environment=HOME="/home/marlise",USER="marlise"
priority=300
depends_on=mod-host

[program:session-manager]
command=python3 -m session_manager.main
directory=/opt/marlise/session-manager
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/session-manager.log
stdout_logfile=/var/log/marlise/session-manager.log
environment=HOME="/home/marlise",USER="marlise",PYTHONPATH="/opt/marlise/session-manager"
priority=400
depends_on=modhost-bridge

[program:client-interface]
command=python3 -m uvicorn api.main:app --host 0.0.0.0 --port 8080 --app-dir /opt/marlise/client-interface/web_api/api
directory=/opt/marlise/client-interface/web_api/api
user=marlise
autostart=true
autorestart=true
stderr_logfile=/var/log/marlise/client-interface.log
stdout_logfile=/var/log/marlise/client-interface.log
environment=HOME="/home/marlise",USER="marlise"
priority=500
depends_on=session-manager