# mado-gov/services.d/browsepy/run
#!/command/with-contenv sh
/usr/local/bin/browsepy --directory /data/user-files --upload /data/user-files --removable /data/user-files 0.0.0.0 8081 \
    2>&1 | sed 's/^/[browsepy] /'